# deploy/avatar-kiosk.service
[Unit]
Description=Algorithmic Absolution avatar kiosk
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/algorithmic-absolution-avatar
ExecStart=/home/pi/algorithmic-absolution-avatar/deploy/start-avatar.sh
Restart=always
RestartSec=5
Environment=DISPLAY=:0
Environment=AVATAR_ORACLE_ONLY=1
# Google TTS key for the proxy (ALG-8). server.js also reads this file directly, so the
# leading "-" makes it optional; keep the key here OR in deploy/tts-proxy/.env, never in git.
EnvironmentFile=-/home/pi/algorithmic-absolution-avatar/deploy/tts-proxy/.env

[Install]
WantedBy=multi-user.target
